# infra/exotalk-zrok.service
# 🧠 EDUCATIONAL CONTEXT: autonomy requires stable endpoints. 
# Using "zrok share reserved" instead of "share public" ensures that the node 
# identity remains consistent across reboots, which is critical for P2P peering.
[Unit]
Description=ExoTalk zrok Signaling Proxy
After=exotalk-signaling.service

[Service]
ExecStart=/home/exocrat/deployments/infra/zrok share reserved exotalkberlin --headless
Restart=always
User=exocrat
WorkingDirectory=/home/exocrat
Environment=HOME=/home/exocrat

[Install]
WantedBy=multi-user.target
