# exo-installer/debian/postinst
#!/bin/sh
set -e
. /usr/share/debconf/confmodule

# 🧠 EDUCATIONAL CONTEXT: The Post-Installation Script
# This script consumes the secure parameters captured by debconf and
# replicates the Rust executor's action sequence to provision the FHS.
# This ensures that a .deb install is IDENTICAL to a TUI install.

if [ "$1" = "configure" ]; then
    # 1. Fetch values from debconf
    db_get conscia-node/components
    COMPONENTS="$RET"
    db_get conscia-node/mesh_namespace
    NAMESPACE="$RET"
    db_get conscia-node/api_port
    PORT="$RET"

    # 2. Create service user (idempotent)
    if ! id -u exo-sys >/dev/null 2>&1; then
        useradd -r -s /usr/sbin/nologin -d /var/lib/exo -m exo-sys
    fi

    # 3. Create FHS directories
    mkdir -p /opt/exo/conscia /opt/exo/ui /etc/exo /var/lib/exo /var/log/exo
    chown -R root:root /opt/exo
    chown root:exo-sys /etc/exo
    chmod 750 /etc/exo
    chown exo-sys:exo-sys /var/lib/exo /var/log/exo
    chmod 700 /var/lib/exo

    # 4. Write config (idempotent)
    cat <<EOF > /etc/exo/conscia.env
# Conscia Daemon Configuration
# Generated by debian/postinst
CONSCIA_MESH_NAMESPACE=$NAMESPACE
CONSCIA_API_PORT=$PORT
CONSCIA_LOG_LEVEL=info
EOF
    chown root:exo-sys /etc/exo/conscia.env
    chmod 640 /etc/exo/conscia.env

    # 5. Reload and Start
    systemctl daemon-reload
    systemctl enable --now exo-conscia.service
fi

# DEBHELPER #

exit 0
